# trinity/deployment/systemd/janus-malaga-briefing.timer
[Unit]
Description=Malaga Briefing Daily Mission (08:00 UTC)
OnFailure=mission-scheduler-failure@malaga_briefing_daily.service

[Timer]
OnCalendar=*-*-* 08:00:00
UTC=true
Persistent=true
Unit=janus-mission-scheduler@malaga_briefing_daily.service

[Install]
WantedBy=timers.target
