# trinity/deployment/systemd/janus-mission-dispatcher.service
[Unit]
Description=Janus Mission Dispatcher Daemon
After=network.target

[Service]
Type=simple
User=balaur
Group=janus
WorkingDirectory=/srv/janus
Environment="PYTHONPATH=/srv/janus:/srv/janus/trinity"
EnvironmentFile=/etc/janus/trinity.env
ExecStart=/usr/bin/python3 /srv/janus/trinity/mission_dispatcher_daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
